# services/full-cycle-btn.service
[Unit]
Description=Full Cycle Button Monitor (GPIO17)
After=default.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__UV_PATH__ run python run_full_cycle_btn.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
