# services/full-cycle-btn.service
[Unit]
Description=Full Cycle Button Monitor (GPIO17)
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/pi/current_state
ExecStart=/home/pi/.local/bin/uv run python run_full_cycle_btn.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
