# services/music-player.service
[Unit]
Description=World Theme Music Player Service
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/pi/daily_mood_theme_song_player
ExecStart=/home/pi/.local/bin/uv run python run_player.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target