# services/process-monitor-web.service
[Unit]
Description=Web Dashboard for AI Music Generation Pipeline
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__UV_PATH__ run uvicorn web.app:app --host 0.0.0.0 --port __WEB_PORT__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
