# services/process-monitor-web.service
[Unit]
Description=Web Dashboard for full process monitor of the project AI music generation project
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/daily_mood_theme_song_player
ExecStart=/home/pi/.local/bin/uv run uvicorn web.app:app --host 127.0.0.1 --port 7070
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target